# files/usr/hlocal/etc/rc.local.d/rc.local
for x in $(ls /etc/rc.local.d/*.rc.conf); do
	. $x
done
